# crontab
*/30 * * * * /usr/local/bin/python3 /usr/src/app/main.py -d >> /var/log/cron.log 2>&1
0 0 1 * * /usr/local/bin/python3 /usr/src/app/main.py -w >> /var/log/cron.log 2>&1
